# modules/shell/apple-silicon/user/.config/profile.d/06-path
# vim: set ft=sh:

export GOPATH="$HOME/.go"
export GEM_HOME=$(ruby -e "puts Gem.user_dir")
export CARGO_HOME="$HOME/.cargo"
export LOCAL_HOME="$HOME/.local"

export PATH="$LOCAL_HOME/bin:$CARGO_HOME/bin:$GOPATH/bin:$GEM_HOME/bin:$PATH"

# Replace macOS coreutils commands with uutils coreutils
export PATH="/opt/homebrew/opt/uutils-coreutils/libexec/uubin:$PATH"
