# modules/shell/common/user/.config/profile.d/00-base
# vim: set ft=sh:

export LANG="en_US.UTF-8"
