# modules/shell/common/user/.config/profile.d/02-defaults
# vim: set ft=sh:

# Default programs
export EDITOR="nvim"
export TERMINAL="wezterm"
export FILE="ranger"

