# modules/shell/common/user/.config/profile.d/06-path
# vim: set ft=sh:

export GOPATH="$HOME/.go"
export GEM_HOME=$(ruby -e "puts Gem.user_dir")

export PATH="$PATH:$HOME/.local/bin:$HOME/.cargo/bin:$GOPATH/bin:$GEM_HOME/bin"
